# install/files/ds-sandbox-openrc
#!/sbin/openrc-run

name="ds-sandbox recycler thing"
command="/bin/ds-sandbox-d"

depend() {
        need networking
}