# install/files/ds-sandbox-openrc
#!/sbin/openrc-run

# http://big-elephants.com/2013-01/writing-your-own-init-scripts/

name="ds-sandbox recycler thing"

depend() {
        need networking
}
start() {
        ebegin "Starting ds-sandbox-d"
        start-stop-daemon --background --start --exec /bin/ds-sandbox-d \
        --make-pidfile --pidfile /root/ds-sandbox-d-pid
        eend $?
}
stop() {
        ebegin "Stopping ds-sandbox-d"
        start-stop-daemon --stop --exec /bin/ds-sandbox-d \
        --pidfile /root/ds-sandbox-d-pid
        eend $?
}
